# recipes-custom/infoip/files/infoip.service
[Unit]
Description=display infoip

[Service]
ExecStart=/usr/bin/infoip.sh

[Install]
WantedBy=multi-user.target
